# data/systemd/fastflx1-andromeda-fs@.service
[Unit]
Description=FastFLX1 Andromeda Permission Guard for %i
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/fastflx1 --monitor andromeda-fs --user %i
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
