# debian/postrm
#!/bin/bash
# Copyright (C) 2026 Alaraajavamma <aki@urheiluaki.fi>
# License: GPL-3.0-or-later

set -e

if [ "$1" = "remove" ]; then
    TARGET_USER=$(id -nu 1000 2>/dev/null || echo "furios")
    TARGET_HOME=$(getent passwd "$TARGET_USER" | cut -d: -f6)

    if [ -d "$TARGET_HOME" ]; then
        echo "Cleaning up Tweak-FLX1s for user $TARGET_USER..."

        if [ -x "/bin/systemctl" ]; then
             echo "Stopping systemd services..."
             systemctl stop tweak-flx1s-andromeda-fs@$TARGET_USER.service 2>/dev/null || true
             systemctl disable tweak-flx1s-andromeda-fs@$TARGET_USER.service 2>/dev/null || true
        fi

        echo "Cleaning up broken systemd user service symlinks..."
        rm -f "$TARGET_HOME/.config/systemd/user/default.target.wants/tweak-flx1s-alarm.service"
        rm -f "$TARGET_HOME/.config/systemd/user/default.target.wants/tweak-flx1s-guard.service"
        rm -f "$TARGET_HOME/.config/systemd/user/default.target.wants/tweak-flx1s-gestures.service"

        echo "Unmounting shared folders..."
        ANDROID_SHARE="$TARGET_HOME/Android-Share"
        LINUX_SHARE="$TARGET_HOME/.local/share/andromeda/data/media/0/Linux-Share"

        umount -l "$ANDROID_SHARE" 2>/dev/null || true
        umount -l "$LINUX_SHARE" 2>/dev/null || true

        if [ -d "$ANDROID_SHARE" ]; then
            echo "Removing Android share mount point..."
            rmdir "$ANDROID_SHARE" 2>/dev/null || true
        fi
        if [ -d "$LINUX_SHARE" ]; then
            echo "Removing Linux share mount point..."
            rmdir "$LINUX_SHARE" 2>/dev/null || true
        fi

        echo "Restoring folder permissions..."
        for dir in "$TARGET_HOME"/*; do
            if [ -d "$dir" ]; then
                dirname=$(basename "$dir")
                if [ "$dirname" != "Android" ]; then
                    echo "Restoring permissions for $dirname..."
                    setfacl -R -b "$dir" 2>/dev/null || true
                    chown -R "$TARGET_USER:$TARGET_USER" "$dir" 2>/dev/null || true
                fi
            fi
        done

        echo "Removing custom components..."
        rm -rf "$TARGET_HOME/.local/share/squeekboard"
        rm -rf "$TARGET_HOME/.local/share/sounds"
        rm -f "$TARGET_HOME/.config/gtk-3.0/gtk.css"

        echo "Removing configuration files..."
        rm -rf "$TARGET_HOME/.config/tweak-flx1s"
        rm -rf "$TARGET_HOME/.cache/tweak-flx1s"

        echo "Tweak-FLX1s cleanup finished."
    fi

    echo "Restoring PAM configuration if backups exist..."
    PAM_FILES=("sudo" "polkit-1" "biomd" "common-password")
    for file in "${PAM_FILES[@]}"; do
        if [ -f "/etc/pam.d/${file}.bak" ]; then
            echo "Restoring /etc/pam.d/${file}..."
            mv "/etc/pam.d/${file}.bak" "/etc/pam.d/${file}"
        fi
    done
fi

#DEBHELPER#
